# deploy/signal_audit/signal-audit-materialize.timer
[Unit]
Description=Refresh signal audit static cards from FMZ JSONL

[Timer]
OnBootSec=30
OnUnitActiveSec=120
Unit=signal-audit-materialize.service

[Install]
WantedBy=timers.target
